# bot.conf
[program:WhisperBot]
directory=/home/ubuntu/WhisperBot
command=python3 app.py

autostart=true
autorestart=true
environment=HOME="/home/ubuntu",USER="ubuntu"

stderr_logfile=/home/ubuntu/WhisperBot/logfile_err.log
stdout_logfile=/home/ubuntu/WhisperBot/logfile.log